# systemd/memo-bridge.service.in
[Unit]
Description=Memo bridge (n8n -> Claude/Codex summarizer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@GROUP@
Environment=HOME=@HOME@
EnvironmentFile=@REPO@/.env
WorkingDirectory=@REPO@/bridge
ExecStart=@PYTHON@ @REPO@/bridge/memo-bridge.py
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
UMask=0077
ReadWritePaths=@REPO@/runtime

[Install]
WantedBy=multi-user.target
