# systemd/n8n-monitor.service.in
[Unit]
Description=Watchdog for n8n, cloudflared, memo bridge and Daily Brief
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=@USER@
Group=@GROUP@
Environment=N8N_ROOT=@REPO@
ExecStart=@PYTHON@ @REPO@/scripts/monitor-system.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=@REPO@/runtime
