# systemd/n8n-monitor.service
[Unit]
Description=Watchdog for n8n, cloudflared, memo bridge and Daily Brief
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=debian
Group=debian
Environment=N8N_ROOT=/SSD/n8n
ExecStart=/usr/bin/python3 /SSD/n8n/scripts/monitor-system.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/SSD/n8n/runtime

